# deploy/monitor-dashboard.service
[Unit]
Description=Options Monitor Dashboard (Streamlit)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=trader
Group=trader
WorkingDirectory=/opt/project_Options
EnvironmentFile=/opt/project_Options/deploy/monitor-dashboard.env
ExecStart=/opt/project_Options/.venv/bin/python -m streamlit run monitor_app.py --server.port 8502
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/opt/project_Options/data /opt/project_Options/logs /opt/project_Options/deploy
UMask=0077

Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target