# deploy/trader-dashboard.service
[Unit]
Description=Weekend Vol Trading Dashboard (Streamlit)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=trader
Group=trader
WorkingDirectory=/opt/project_Options
ExecStart=/opt/project_Options/.venv/bin/python -m streamlit run trader/dashboard.py -- --config configs/trader/weekend_vol_btc.yaml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
